# roles/zfs-exporter/templates/zfs-exporter.service.j2
[Unit]
Description=ZFS exporter

[Service]
Type=simple
DynamicUser=yes
User=zfs_exporter
ExecStart={{ zfs_exporter_bin }}

[Install]
WantedBy=multi-user.target
